# inbox/launchd/com.m13v.fazm-founder-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m13v.fazm-founder-chat</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/matthewdi/fazm/inbox/skill/check-founder-chat.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>15</integer>
	<key>StandardOutPath</key>
	<string>/Users/matthewdi/fazm/inbox/skill/logs/founder-chat-launchd-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/matthewdi/fazm/inbox/skill/logs/founder-chat-launchd-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/matthewdi/.nvm/versions/node/v20.19.4/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/matthewdi</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
